# reports/sigma_rules/t1547_001.yml
title: Registry Run Key Persistence (GAN-Generated)
id: 4664336a-2343-57c9-927b-c19ca7948984
status: experimental
description: |
    Auto-generated by find-evil-2026 GAN adversarial training loop.
    Trained iterations: 3000
    Detection weight: 50
    Evasions observed during training: 148
    Adversarial validation: NEEDS_HARDENING (148 evasions tested, 131 bypassed)
    Source: GAN_trained
references:
    - https://attack.mitre.org/techniques/T1547/001/
    - https://github.com/OTRF/Security-Datasets (Mordor training data)
author: find-evil-2026-gan
date: 2026/05/14
modified: 2026/05/14
tags:
    - attack.persistence
    - attack.t1547.001
logsource:
    product: windows
    category: registry_set
detection:
    sel_registry:
        TargetObject|contains:
            - '\lsass'
            - 'Set-ItemProperty HKCU\Software\Microsoft\Windows\CurrentVersion\Run'
    sel_keywords:
        EventData|contains:
            - 'dllhost\\svchost'
            - psexesvc
            - WmiPrvSE.exe
            - 'Network conn allowed via kernel filter layer SourcePort:445'
    condition: sel_registry or sel_keywords
falsepositives:
    - Legitimate software installations and Windows updates
    - Adversarially validated against 148 Red Agent evasions
level: high
fields:
    - Image
    - CommandLine
    - TargetObject
    - TargetImage
    - DestinationIp
# GAN metadata (machine-readable)
# gan_iteration: 3000
# bypass_rate: 0.885
# adversarially_validated: false
